# scripts/oracle_systemd.service
[Unit]
Description=StrategyHarvester 24/7 Signal Engine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# __USER__ and __PROJECT_DIR__ are substituted by setup_oracle.sh.
User=__USER__
WorkingDirectory=__PROJECT_DIR__
ExecStart=__PROJECT_DIR__/.venv/bin/python -m scheduler.runner_prod
Restart=always
RestartSec=10
# Capture stdout/stderr into journald (view with: journalctl -u stratharv -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stratharv

[Install]
WantedBy=multi-user.target
