# scripts/oracle_dashboard.service
[Unit]
Description=StrategyHarvester live dashboard (management UI, Phase 3)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# __USER__ and __PROJECT_DIR__ are substituted by scripts/setup_oracle_timers.sh.
User=__USER__
WorkingDirectory=__PROJECT_DIR__
ExecStart=__PROJECT_DIR__/.venv/bin/python -m uvicorn monitoring.dashboard_server:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stratharv-dashboard

[Install]
WantedBy=multi-user.target
