# scripts/oracle_content_intel.service
[Unit]
Description=StrategyHarvester Content Intelligence watchlist (one pass)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# __USER__ and __PROJECT_DIR__ are substituted by setup_oracle.sh / manual install.
User=__USER__
WorkingDirectory=__PROJECT_DIR__
ExecStart=__PROJECT_DIR__/.venv/bin/python -m scheduler.content_intelligence_cron
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stratharv-content-intel
